# lib/go/liora-ai.service
[Unit]
Description=Liora AI Service
Documentation=https://github.com/your-repo/lib/go
After=network-online.target postgresql.service redis.service
Wants=network-online.target
Requires=postgresql.service

[Service]
Type=simple
User=liora
Group=liora
WorkingDirectory=/root/liora/lib/go
Environment="PATH=/usr/local/go/bin:/usr/bin:/bin"
Environment="GOMAXPROCS=4"

# Service binary
ExecStart=/root/liora/lib/go/bin/ai-service

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=5

# Resource limits
LimitNOFILE=65536
MemoryMax=2G
CPUQuota=200%

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/root/liora/lib/go/logs

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=liora-ai

# Graceful shutdown
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target